# systemd/sriov-vf-mappings.service
[Unit]
Description=Sync Proxmox PCI resource mappings for SR-IOV VFs
# Needs the VFs to exist (sriov-vfs.service) AND pmxcfs to be mounted
# (pve-cluster.service, which is ordered After=network.target). Must finish
# before guests that reference the mappings are started.
Requires=sriov-vfs.service
After=sriov-vfs.service pve-cluster.service
Wants=pve-cluster.service
Before=pve-guests.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/opt/schlarpc/bin/sync-sriov-vf-mappings

[Install]
WantedBy=multi-user.target
