# systemd/sriov-vf-mappings@.service
[Unit]
Description=Sync Proxmox PCI resource mappings for SR-IOV VFs on %i
# Needs the VFs on this PF to exist (sriov-vfs@%i) AND pmxcfs to be mounted
# (pve-cluster.service, ordered After=network.target). Must finish before
# guests that reference the mappings are started.
Requires=sriov-vfs@%i.service
After=sriov-vfs@%i.service pve-cluster.service
Wants=pve-cluster.service
Before=pve-guests.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/opt/schlarpc/bin/sync-sriov-vf-mappings %i

# No [Install] section: this is pulled in by sriov-vfs@%i.service (Wants=), so
# you enable only that one. `systemctl enable sriov-vf-mappings@<pf>` is neither
# needed nor possible.
