# systemd/sriov-vfs.service
[Unit]
Description=Create Mellanox SR-IOV VFs (switchdev + per-VF MAC)
# The eswitch representors created here are enslaved into a vlan-aware bridge
# in /etc/network/interfaces, so this must finish before networking is brought
# up. It deliberately has NO pve-* dependency -- it no longer calls pvesh.
Before=network-pre.target network-online.target
Wants=network-pre.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/opt/schlarpc/bin/create-sriov-vfs

[Install]
WantedBy=multi-user.target
