# debian/postinst
#!/bin/sh
set -e

if [ "$1" = configure ]; then
	if [ -d /run/systemd/system ]; then
		systemctl --system daemon-reload >/dev/null 2>&1 || true
	fi
	# Enable for next boot. Deliberately NOT started now: these are boot-time
	# provisioning units, and re-running VF creation under live guests is
	# undesirable. They take effect on the next reboot.
	if command -v deb-systemd-helper >/dev/null 2>&1; then
		deb-systemd-helper enable sriov-vfs.service >/dev/null 2>&1 || true
		deb-systemd-helper enable sriov-vf-mappings.service >/dev/null 2>&1 || true
	else
		systemctl enable sriov-vfs.service sriov-vf-mappings.service >/dev/null 2>&1 || true
	fi
fi

exit 0
