# debian/postinst
#!/bin/sh
set -e

if [ "$1" = configure ]; then
	if [ -d /run/systemd/system ]; then
		systemctl --system daemon-reload >/dev/null 2>&1 || true
	fi
	# These are template units, so nothing is auto-enabled: the package cannot
	# know which PF(s) you provision. Enable one unit per PF; the matching
	# mapping sync is pulled in automatically:
	#   systemctl enable mlx5-sriov-vfs@enp33s0f0np0.service
fi

exit 0
