# debian/prerm
#!/bin/sh
set -e

if [ "$1" = remove ] || [ "$1" = deconfigure ]; then
	if command -v deb-systemd-helper >/dev/null 2>&1; then
		deb-systemd-helper --no-reload disable sriov-vfs.service >/dev/null 2>&1 || true
		deb-systemd-helper --no-reload disable sriov-vf-mappings.service >/dev/null 2>&1 || true
	fi
	# VFs already created this boot are left in place; removing the package does
	# not tear down a running datapath.
fi

exit 0
